# docker/run/fs/etc/supervisor/conf.d/vnc.conf
# ==============================================================================
# VNC / X11 stack for Agent Zero remote desktop
#
# These services are non-essential — the Web UI (port 80) works without them.
# The settings below are tuned to:
#   1. Clean up stale X lock files on (re)start  (start_xvfb.sh wrapper)
#   2. Tolerate the race between Xvfb and its dependents (generous retries)
#   3. Never enter FATAL state, which would trigger the event listener
#      to kill supervisord and tear down the entire container
# ==============================================================================

[program:setup_vnc_password]
command=/exe/setup_vnc_password.sh
autostart=true
autorestart=false
priority=5
startsecs=0
exitcodes=0
stderr_logfile=/var/log/supervisor/setup_vnc_password.err.log
stdout_logfile=/var/log/supervisor/setup_vnc_password.out.log

[program:xvfb]
# Reason: The wrapper script removes /tmp/.X99-lock and /tmp/.X11-unix/X99
# before exec-ing Xvfb, preventing the "Server is already active" crash loop.
command=/exe/start_xvfb.sh
autostart=true
autorestart=true
priority=10
startsecs=2
startretries=5
stderr_logfile=/var/log/supervisor/xvfb.err.log
stdout_logfile=/var/log/supervisor/xvfb.out.log
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB

[program:setup_clipboard_shortcuts]
command=/usr/local/bin/setup-clipboard-shortcuts
environment=DISPLAY=":99"
autostart=true
autorestart=false
priority=15
startsecs=0
exitcodes=0
stderr_logfile=/var/log/supervisor/setup_clipboard.err.log
stdout_logfile=/var/log/supervisor/setup_clipboard.out.log

[program:fluxbox]
command=/usr/bin/fluxbox
environment=DISPLAY=":99"
autostart=true
autorestart=true
priority=20
# Reason: fluxbox needs Xvfb to be fully up.  startsecs=0 means supervisor
# will not mark it FATAL on early exit; startretries gives it time to wait
# for display :99 to become available.
startsecs=0
startretries=10
stderr_logfile=/var/log/supervisor/fluxbox.err.log
stdout_logfile=/var/log/supervisor/fluxbox.out.log
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB

[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -rfbauth /root/.vnc/passwd -listen 0.0.0.0 -xkb -forever -shared -rfbport 5900
environment=DISPLAY=":99"
autostart=true
autorestart=true
priority=30
# Reason: x11vnc needs display :99.  Same resilience pattern as fluxbox.
startsecs=0
startretries=10
stderr_logfile=/var/log/supervisor/x11vnc.err.log
stdout_logfile=/var/log/supervisor/x11vnc.out.log
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB

[program:autocutsel]
# Reason: Removed the -fork flag.  With -fork, autocutsel daemonizes and
# the parent exits immediately — supervisor interprets this as a crash.
# Without -fork, autocutsel stays in the foreground so supervisor can
# manage its lifecycle correctly.
command=/usr/bin/autocutsel -selection CLIPBOARD
environment=DISPLAY=":99"
autostart=true
autorestart=true
priority=40
startsecs=0
startretries=10
stderr_logfile=/var/log/supervisor/autocutsel.err.log
stdout_logfile=/var/log/supervisor/autocutsel.out.log
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB

[program:wallpaper]
command=/bin/bash -c "sleep 5 && feh --bg-fill /usr/share/wallpapers/agent-zero.jpg"
environment=DISPLAY=":99"
autostart=true
autorestart=false
priority=50
startsecs=0
exitcodes=0
stderr_logfile=/var/log/supervisor/wallpaper.err.log
stdout_logfile=/var/log/supervisor/wallpaper.out.log
