# docker/run/fs/etc/supervisor/conf.d/vnc.conf
[program:setup_vnc_password]
command=/exe/setup_vnc_password.sh
autostart=true
autorestart=false
priority=5
startsecs=0
exitcodes=0
stderr_logfile=/var/log/supervisor/setup_vnc_password.err.log
stdout_logfile=/var/log/supervisor/setup_vnc_password.out.log

[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1920x1080x24 -ac +extension GLX +render -noreset
autostart=true
autorestart=true
priority=10
stderr_logfile=/var/log/supervisor/xvfb.err.log
stdout_logfile=/var/log/supervisor/xvfb.out.log

[program:setup_clipboard_shortcuts]
command=/usr/local/bin/setup-clipboard-shortcuts
environment=DISPLAY=":99"
autostart=true
autorestart=false
priority=15
startsecs=0
exitcodes=0
stderr_logfile=/var/log/supervisor/setup_clipboard.err.log
stdout_logfile=/var/log/supervisor/setup_clipboard.out.log

[program:fluxbox]
command=/usr/bin/fluxbox
environment=DISPLAY=":99"
autostart=true
autorestart=true
priority=20
stderr_logfile=/var/log/supervisor/fluxbox.err.log
stdout_logfile=/var/log/supervisor/fluxbox.out.log

[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -passwd vnc123 -listen 0.0.0.0 -xkb -forever -shared -rfbport 5900
environment=DISPLAY=":99"
autostart=true
autorestart=true
priority=30
stderr_logfile=/var/log/supervisor/x11vnc.err.log
stdout_logfile=/var/log/supervisor/x11vnc.out.log

[program:autocutsel]
command=/usr/bin/autocutsel -fork -selection PRIMARY -selection CLIPBOARD
environment=DISPLAY=":99"
autostart=true
autorestart=true
priority=40
stderr_logfile=/var/log/supervisor/autocutsel.err.log
stdout_logfile=/var/log/supervisor/autocutsel.out.log

[program:wallpaper]
command=/bin/bash -c "sleep 5 && feh --bg-fill /usr/share/wallpapers/agent-zero.jpg"
environment=DISPLAY=":99"
autostart=true
autorestart=false
priority=50
startsecs=0
exitcodes=0
stderr_logfile=/var/log/supervisor/wallpaper.err.log
stdout_logfile=/var/log/supervisor/wallpaper.out.log
